# MinecraftServer.service
[Unit]
Description=minecraft server

[Service]
ExecStart=/bin/bash ~/minecraft/casually-chaotic/start.sh

[Install]
WantedBy=multi-user.target