# MinecraftServer.service
[Unit]
Description=Minecraft Server

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=~/minecraft/casually-chaotic/start.sh
Restart=always