# MinecraftServer.socket
[Unit]
PartOf=MinecraftServer.service

[Socket]
ListenFIFO=%t/MinecraftServer.stdin